# scripts/mqtt-master.service
[Unit]
Description=MQTT Master
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=mqtt-master
Group=mqtt-master
WorkingDirectory=/opt/mqtt-master
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
